# qtcm/quectel-CM/files/etc/init.d/qtcm
#!/bin/sh /etc/rc.common

START=99
USE_PROCD=1

. /lib/functions.sh

CONF="qtcm"
SECTION="main"
EXTRA_COMMANDS="killpdn"
EXTRA_HELP="    killpdn [pdp]    Disconnect a specific PDN using quectel-CM"

append_arg() {
	local value="$1"
	local flag="$2"

	[ -n "$value" ] && procd_append_param command "$flag" "$value"
}

load_settings() {
	config_load "$CONF"

	config_get enabled "$SECTION" enabled "1"
	config_get log_enabled "$SECTION" log "1"
	config_get log_file "$SECTION" log_file "/tmp/q-cm.log"
	config_get usbmon_log_file "$SECTION" usbmon_log_file ""
	config_get cell_internet_mode "$SECTION" cell_internet_mode "nat"
	config_get network_interface "$SECTION" network_interface ""
	config_get pincode "$SECTION" pincode ""
	config_get proxy_mode "$SECTION" proxy_mode ""
	config_get pdp "$SECTION" pdp "1"
	config_get mux_id "$SECTION" mux_id ""
	config_get ip_type "$SECTION" ip_type "ipv4"
	config_get verbose "$SECTION" verbose "0"
	config_get no_dhcp "$SECTION" no_dhcp "0"
	config_get auto_apn "$SECTION" auto_apn "1"
	config_get apn "$SECTION" apn ""
	config_get username "$SECTION" username ""
	config_get password "$SECTION" password ""
	config_get auth "$SECTION" auth "none"

	bridge_mode="0"
	[ "$cell_internet_mode" = "ippt" ] && bridge_mode="1"
	[ "$auto_apn" = "1" ] && apn=""
}

append_profile_args() {
	[ "$auto_apn" = "1" ] && return 0
	[ -z "$apn" ] && return 0

	procd_append_param command -s "$apn"

	# quectel-CM parses user/password/auth positionally after -s.
	# Only append trailing credentials when both user and password exist,
	# otherwise later values would shift into the wrong argument slot.
	if [ -n "$username" ] && [ -n "$password" ]; then
		procd_append_param command "$username" "$password"

		case "$auth" in
			pap) procd_append_param command 1 ;;
			chap) procd_append_param command 2 ;;
			mschapv2) procd_append_param command 3 ;;
			none) procd_append_param command 0 ;;
		esac
	fi
}

ensure_qmi_wwan_q_loaded() {
	[ -d /sys/module/qmi_wwan_q ] && return 0

	insmod qmi_wwan_q.ko qmap_mode=1 2>/dev/null
}

boot_prepare() {
	rmmod usbserial 2>/dev/null
	rmmod ch341 2>/dev/null
}

boot_finish() {
	modprobe usbserial 2>/dev/null
	sleep 2
	modprobe ch341 2>/dev/null
}

start_service() {
	load_settings
	[ "$enabled" = "0" ] && return 0

	echo "start wait 30sec >>>>>>>>"
	boot_prepare
	sleep 30
	echo "Module loading >>>>>>>>"
	ensure_qmi_wwan_q_loaded
	sleep 5
	echo "############## Start LTE ############"

	procd_open_instance
	procd_set_param command /sbin/quectel-CM

	[ "$bridge_mode" = "1" ] && procd_append_param command -b
	append_arg "$pincode" -p
	append_arg "$proxy_mode" -p
	append_arg "$network_interface" -i
	[ "$pdp" != "1" ] && append_arg "$pdp" -n
	[ -n "$mux_id" ] && append_arg "$mux_id" -m
	[ "$ip_type" = "ipv4v6" ] && procd_append_param command -4 -6
	[ "$ip_type" = "ipv6" ] && procd_append_param command -6
	append_profile_args
	[ "$verbose" = "1" ] && procd_append_param command -v
	[ "$no_dhcp" = "1" ] && procd_append_param command -d
	[ "$log_enabled" = "1" ] && append_arg "$log_file" -f
	append_arg "$usbmon_log_file" -u

	procd_set_param respawn
	procd_close_instance

	echo "####################################"
	boot_finish
}

stop_service() {
	service_stop /sbin/quectel-CM
}

killpdn() {
	local pdp="${1:-1}"

	/sbin/quectel-CM -k "$pdp"
}

service_triggers() {
	procd_add_reload_trigger "$CONF"
}
